# cytech_comfort_mqtt/rootfs/etc/services.d/cytech_comfort_mqtt/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Comfort to MQTT..."
exec python /comfort2/bridge.py